# etc/supervisord.d/zeppelin.ini
[program:zeppelin]
command = /opt/zeppelin/bin/zeppelin.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
